# crontab
0,5,10,20,25,30,35,40,45,50 * * * * /usr/bin/python /home/pi/SCIAproject/piCameraControl.py >/dev/null 2>&1
